# share/discord-wrangler.service.in
[Unit]
Description=discord-wrangler NFQUEUE daemon
Documentation=file://@DOCDIR@/README.md
After=nftables.service network.target
Wants=nftables.service

[Service]
Type=simple
User=discord-wrangler
Group=discord-wrangler
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateDevices=true
PrivateTmp=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
SystemCallFilter=@system-service ~@privileged
Environment=WRANGLER_QUEUE_NUM=@QUEUE_NUM@
Environment=WRANGLER_LOG_LEVEL=info
ReadOnlyPaths=/sys/fs/cgroup
ExecStartPre=/usr/sbin/nft -f /etc/nftables.d/discord-wrangler.nft
ExecStart=@SBINDIR@/discord-wranglerd
ExecStopPost=-/usr/sbin/nft delete table inet discord_wrangler
ExecStopPost=-@SBINDIR@/discord-wrangler-cleanup
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
